# deployment/staging/config/supervisord.staging.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info
silent=false

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ============================================
# Service Groups for Easy Management
# ============================================
[group:core_services]
programs=auth_service,account_service,session_service,authorization_service
priority=10

[group:business_services]
programs=payment_service,wallet_service,order_service,organization_service,invitation_service
priority=20

[group:storage_services]
programs=storage_service,album_service,media_service,vault_service
priority=30

[group:device_services]
programs=device_service,ota_service,telemetry_service
priority=40

[group:support_services]
programs=audit_service,notification_service,task_service,event_service,product_service,billing_service,calendar_service,weather_service,memory_service
priority=50

# ============================================
# Core Services
# ============================================
[program:auth_service]
command=uvicorn microservices.auth_service.main:app --host 0.0.0.0 --port 8201 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="auth_service"

[program:account_service]
command=uvicorn microservices.account_service.main:app --host 0.0.0.0 --port 8202 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="account_service"

[program:session_service]
command=uvicorn microservices.session_service.main:app --host 0.0.0.0 --port 8203 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="session_service"

[program:authorization_service]
command=uvicorn microservices.authorization_service.main:app --host 0.0.0.0 --port 8204 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="authorization_service"

# ============================================
# Business Services
# ============================================
[program:payment_service]
command=uvicorn microservices.payment_service.main:app --host 0.0.0.0 --port 8207 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="payment_service"

[program:wallet_service]
command=uvicorn microservices.wallet_service.main:app --host 0.0.0.0 --port 8208 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="wallet_service"

[program:order_service]
command=uvicorn microservices.order_service.main:app --host 0.0.0.0 --port 8210 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="order_service"

[program:organization_service]
command=uvicorn microservices.organization_service.main:app --host 0.0.0.0 --port 8212 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="organization_service"

[program:invitation_service]
command=uvicorn microservices.invitation_service.main:app --host 0.0.0.0 --port 8213 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="invitation_service"

# ============================================
# Storage Services
# ============================================
[program:storage_service]
command=uvicorn microservices.storage_service.main:app --host 0.0.0.0 --port 8209 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="storage_service"

[program:album_service]
command=uvicorn microservices.album_service.main:app --host 0.0.0.0 --port 8219 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="album_service"

[program:media_service]
command=uvicorn microservices.media_service.main:app --host 0.0.0.0 --port 8222 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="media_service"

[program:vault_service]
command=uvicorn microservices.vault_service.main:app --host 0.0.0.0 --port 8214 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="vault_service"

# ============================================
# Device Services
# ============================================
[program:device_service]
command=uvicorn microservices.device_service.main:app --host 0.0.0.0 --port 8220 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="device_service"

[program:ota_service]
command=uvicorn microservices.ota_service.main:app --host 0.0.0.0 --port 8221 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="ota_service"

[program:telemetry_service]
command=uvicorn microservices.telemetry_service.main:app --host 0.0.0.0 --port 8225 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="telemetry_service"

# ============================================
# Support Services
# ============================================
[program:audit_service]
command=uvicorn microservices.audit_service.main:app --host 0.0.0.0 --port 8205 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="audit_service"

[program:notification_service]
command=uvicorn microservices.notification_service.main:app --host 0.0.0.0 --port 8206 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="notification_service"

[program:task_service]
command=uvicorn microservices.task_service.main:app --host 0.0.0.0 --port 8211 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="task_service"

[program:event_service]
command=uvicorn microservices.event_service.main:app --host 0.0.0.0 --port 8230 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="event_service"

[program:product_service]
command=uvicorn microservices.product_service.main:app --host 0.0.0.0 --port 8215 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="product_service"

[program:billing_service]
command=uvicorn microservices.billing_service.main:app --host 0.0.0.0 --port 8216 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="billing_service"

[program:calendar_service]
command=uvicorn microservices.calendar_service.main:app --host 0.0.0.0 --port 8217 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="calendar_service"

[program:weather_service]
command=uvicorn microservices.weather_service.main:app --host 0.0.0.0 --port 8218 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="weather_service"

[program:memory_service]
command=uvicorn microservices.memory_service.main:app --host 0.0.0.0 --port 8223 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="/app",SERVICE_NAME="memory_service"
