# deployment/staging/config/supervisord.staging.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info
silent=false

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ============================================
# Core Services (Priority 10)
# ============================================
[program:auth_service]
command=python -u -m uvicorn microservices.auth_service.main:app --host 0.0.0.0 --port 8201 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="auth_service"

[program:account_service]
command=python -u -m uvicorn microservices.account_service.main:app --host 0.0.0.0 --port 8202 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="account_service"

[program:session_service]
command=python -u -m uvicorn microservices.session_service.main:app --host 0.0.0.0 --port 8203 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="session_service"

[program:authorization_service]
command=python -u -m uvicorn microservices.authorization_service.main:app --host 0.0.0.0 --port 8204 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="authorization_service"

# ============================================
# Business Services (Priority 20)
# ============================================
[program:payment_service]
command=python -u -m uvicorn microservices.payment_service.main:app --host 0.0.0.0 --port 8207 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="payment_service"

[program:wallet_service]
command=python -u -m uvicorn microservices.wallet_service.main:app --host 0.0.0.0 --port 8208 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="wallet_service"

[program:order_service]
command=python -u -m uvicorn microservices.order_service.main:app --host 0.0.0.0 --port 8210 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="order_service"

[program:organization_service]
command=python -u -m uvicorn microservices.organization_service.main:app --host 0.0.0.0 --port 8212 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="organization_service"

[program:invitation_service]
command=python -u -m uvicorn microservices.invitation_service.main:app --host 0.0.0.0 --port 8213 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="invitation_service"

# ============================================
# Storage Services (Priority 30)
# ============================================
[program:storage_service]
command=python -u -m uvicorn microservices.storage_service.main:app --host 0.0.0.0 --port 8209 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="storage_service",STORAGE_SERVICE_PORT="8209"

[program:album_service]
command=python -u -m uvicorn microservices.album_service.main:app --host 0.0.0.0 --port 8219 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="album_service"

[program:media_service]
command=python -u -m uvicorn microservices.media_service.main:app --host 0.0.0.0 --port 8222 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="media_service"

[program:vault_service]
command=python -u -m uvicorn microservices.vault_service.main:app --host 0.0.0.0 --port 8214 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="vault_service"

# ============================================
# Device Services (Priority 40)
# ============================================
[program:device_service]
command=python -u -m uvicorn microservices.device_service.main:app --host 0.0.0.0 --port 8220 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="device_service"

[program:ota_service]
command=python -u -m uvicorn microservices.ota_service.main:app --host 0.0.0.0 --port 8221 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="ota_service"

[program:telemetry_service]
command=python -u -m uvicorn microservices.telemetry_service.main:app --host 0.0.0.0 --port 8225 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="telemetry_service"

[program:location_service]
command=python -u -m uvicorn microservices.location_service.main:app --host 0.0.0.0 --port 8224 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="location_service"

# ============================================
# Support Services (Priority 50)
# ============================================
[program:audit_service]
command=python -u -m uvicorn microservices.audit_service.main:app --host 0.0.0.0 --port 8205 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="audit_service"

[program:notification_service]
command=python -u -m uvicorn microservices.notification_service.main:app --host 0.0.0.0 --port 8206 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="notification_service"

[program:task_service]
command=python -u -m uvicorn microservices.task_service.main:app --host 0.0.0.0 --port 8211 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="task_service"

[program:event_service]
command=python -u -m uvicorn microservices.event_service.main:app --host 0.0.0.0 --port 8230 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="event_service"

[program:product_service]
command=python -u -m uvicorn microservices.product_service.main:app --host 0.0.0.0 --port 8215 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="product_service"

[program:billing_service]
command=python -u -m uvicorn microservices.billing_service.main:app --host 0.0.0.0 --port 8216 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="billing_service"

[program:calendar_service]
command=python -u -m uvicorn microservices.calendar_service.main:app --host 0.0.0.0 --port 8217 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="calendar_service"

[program:weather_service]
command=python -u -m uvicorn microservices.weather_service.main:app --host 0.0.0.0 --port 8218 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="weather_service"

[program:memory_service]
command=python -u -m uvicorn microservices.memory_service.main:app --host 0.0.0.0 --port 8223 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="memory_service"

[program:compliance_service]
command=python -u -m uvicorn microservices.compliance_service.main:app --host 0.0.0.0 --port 8226 --reload --log-level info --log-config /app/deployment/staging/config/logging.json
directory=/app
user=app
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONPATH="/app",SERVICE_NAME="compliance_service"
